# deploy/supervisord.conf
; CodeAgent Multi-Agent Supervisor config
; Install: supervisorctl reread && supervisorctl update

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/home/u2/.codeagent/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
user=u2

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:codeagent]
command=/home/u2/agent/env/bin/codeagent serve --port 8088
directory=/home/u2/agent/codeagent
user=u2
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10
stdout_logfile=/home/u2/.codeagent/logs/codeagent-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/home/u2/.codeagent/logs/codeagent-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED=1
